# migration/deploy/smelitehajj-backend.service
[Unit]
Description=SM Elite Hajj Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/smelitehajj/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=smelitehajj-backend
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
